# scripts/systemd/allura-memory-stack.service.template
[Unit]
Description=Allura Memory stack
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=__REPO_ROOT__
ExecStart=/usr/bin/env bash -lc 'cd "__REPO_ROOT__" && ./scripts/brain-stack.sh up'
ExecStop=/usr/bin/env bash -lc 'cd "__REPO_ROOT__" && ./scripts/brain-stack.sh down'
TimeoutStartSec=180
TimeoutStopSec=120

[Install]
WantedBy=default.target
